# extern/pyre/bot/slaves/beryllium/pyre.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnvironmentVariables</key>
	<dict>
		<key>LANG</key>
		<string>en_US.UTF-8</string>
		<key>PATH</key>
		<string>/opt/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
	</dict>
	<key>Label</key>
	<string>buildslave</string>
	<key>ProgramArguments</key>
	<array>
		<string>/opt/local/bin/buildslave</string>
		<string>start</string>
		<string>/Users/mga/buildbot/osx</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
</dict>
</plist>
